# docker/s6/ntpd/run
#!/bin/bash

[ -f setup ] && . setup

exec gosu root /usr/sbin/ntpd -d
