# docker/s6/ntpd/run
#!/bin/bash

[ -f setup ] && . setup

exec gosu root /usr/sbin/ntpd -n -b -g -l stdout -u ntp:ntp
