# docker/s6/.s6-svscan/finish
#!/bin/bash

# Cleanup SOCAT services and s6 event folder
rm -rf $(find /app/s6/ -name 'event')
rm -rf /app/s6/SOCAT_*
