# docker/s6/cron/run
#!/bin/bash

[ -f setup ] && . setup

exec cron -f
